# homework2/dx-scraper.service
[Unit]
Description=DX Cluster Scraper Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=steve
Group=steve
WorkingDirectory=/home/steve/GITHUB/cs330-projects/homework2
Environment=PATH=/home/steve/GITHUB/cs330-projects/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/steve/GITHUB/cs330-projects/.venv/bin/python /home/steve/GITHUB/cs330-projects/homework2/dx_cluster_live_pg.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dx-scraper

# Environment file
EnvironmentFile=/home/steve/GITHUB/cs330-projects/homework2/.env

[Install]
WantedBy=multi-user.target